# annoying-bot.service
# SPDX-FileCopyrightText: 2026 Peter Lemenkov <lemenkov@gmail.com>
# SPDX-License-Identifier: MIT

[Unit]
Description=Annoying Men's Phrases Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
DynamicUser=yes
WorkingDirectory=/opt/annoying_bot
ExecStart=/usr/bin/python3 bot.py
Restart=on-failure
RestartSec=10

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/annoying_bot

[Install]
WantedBy=multi-user.target
